# templates/activemq.conf.erb
start on startup
respawn
script
  exec /usr/ActiveMQ/bin/activemq console
end script